# docker-vertica/s6-rc.d/cron/run
#!/command/execlineb -P
sudo /usr/sbin/cron -f
