# docker-vertica/s6-rc.d/cron/run
#!/command/execlineb -P
/usr/sbin/cron -f
