# docker-vertica-v2/s6-rc.d/nma/run
#!/command/execlineb -P
# explicitly load container environment
with-contenv
# start node_management_agent
/opt/vertica/bin/node_management_agent
