# docker-vertica-v2/s6-rc.d/nma/up
#!/command/execlineb -P
/opt/vertica/bin/manage_node_agent.sh start
