# docker-vertica-v2/s6-rc.d/nma/up
#!/command/execlineb -P
# explicitly load container environment
with-contenv
# start node_management_agent
/opt/vertica/bin/manage_node_agent.sh start
