# docker-vertica/s6-rc.d/sshd/run
#!/command/execlineb -P
foreground { sudo ssh-keygen -q -A } sudo /usr/sbin/sshd


